# launchd/com.remarkable.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.remarkable.weekly</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>0</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/jim/dev/remarkable_ingest/.venv/bin/python</string>
      <string>/Users/jim/dev/remarkable_ingest/weekly_summary.py</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/jim/dev/remarkable_ingest/weekly_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jim/dev/remarkable_ingest/weekly_stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/jim/dev/remarkable_ingest</string>
  </dict>
</plist>

