# launchd/com.remarkable.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key><string>com.remarkable.ingest</string>
    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>RunAtLoad</key><true/>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/jim/dev/remarkable_ingest/.venv/bin/python</string>
      <string>/Users/jim/dev/remarkable_ingest/main.py</string>
    </array>
    <key>StandardOutPath</key><string>/Users/jim/dev/remarkable_ingest/stdout.log</string>
    <key>StandardErrorPath</key><string>/Users/jim/dev/remarkable_ingest/stderr.log</string>
    <key>WorkingDirectory</key><string>/Users/jim/dev/remarkable_ingest</string>
  </dict>
</plist>

